# roles/certbot/templates/cert-expiry-metric.service.j2
[Unit]
Description=Write cert expiry textfile metric for node_exporter
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/letsencrypt/renewal-hooks/deploy/cert-expiry-metric.sh
